# systemd/ns4j.service
[Unit]
Description=ns4j — Azure DDNS update (oneshot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/node /opt/ns4j/index.mjs
EnvironmentFile=/etc/ns4j/.env
User=ns4j
Group=ns4j

# State directory on tmpfs
RuntimeDirectory=ns4j
RuntimeDirectoryMode=0750

# Security hardening
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
